# lib/systemd/system/wpa_supplicant.timer
[Unit]
Description=Run wpa_supplicant boot

[Timer]
OnBootSec=2min

[Install]
WantedBy=timers.target
